# pi/systemd/generic.service
[Unit]
Description=Generic Systemd Service File

[Service]

